# launchd/com.jmclaw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jmclaw</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/node@22/bin/node</string>
        <string>--require</string>
        <string>/Users/xiangyuanmeng/Documents/clawpool/jmclaw/proxy-preload.cjs</string>
        <string>/Users/xiangyuanmeng/Documents/clawpool/jmclaw/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/xiangyuanmeng/Documents/clawpool/jmclaw</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/opt/node@22/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/xiangyuanmeng</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/xiangyuanmeng/Documents/clawpool/jmclaw/nanoclaw.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/xiangyuanmeng/Documents/clawpool/jmclaw/nanoclaw.log</string>
</dict>
</plist>
